# luci-theme-netgear/root/etc/uci-defaults/30_luci-theme-netgear
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Netgear=/luci-static/netgear
	commit luci
EOF
exit 0 
# ↑↑↑ DO NOT EDIT THIS LINE
